# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

# =============================================================================
# Graceful Shutdown Configuration
# Requirements: 13.4 - Support graceful shutdown with configurable timeout periods
# =============================================================================
# When supervisord receives SIGTERM, it will stop all managed processes
# in reverse priority order (highest priority number first)

# =============================================================================
# PocketBase Backend Service
# Priority 100 - Starts first as other services depend on it
# Stops last during shutdown to ensure database is available
# Requirements: 1.3, 1.4
# =============================================================================
[program:pocketbase]
command=/app/pb/pocketbase serve --http=0.0.0.0:8090 --dir=%(ENV_PB_DATA_DIR)s --migrationsDir=/app/pb/pb_migrations --hooksDir=/app/pb/pb_hooks
directory=/app/pb
autostart=true
autorestart=true
# Log to console for visibility
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_backups=0
stdout_logfile_backups=0
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
user=nextjs
environment=HOME="/app",USER="nextjs"
# Start first with highest priority (lower number = higher priority)
priority=100
# Consider started after 5 seconds of running (give time for initialization)
startsecs=5
startretries=5
# Graceful shutdown: send SIGTERM, wait for graceful shutdown
stopsignal=TERM
stopwaitsecs=%(ENV_GRACEFUL_SHUTDOWN_TIMEOUT)s
# Kill with SIGKILL if SIGTERM doesn't work
killasgroup=true
stopasgroup=true

# =============================================================================
# Next.js Webapp Service
# Priority 200 - Starts after PocketBase
# =============================================================================
[program:nextjs]
command=yarn workspace @project/webapp start
directory=/app
autostart=true
autorestart=true
# Log to console for visibility
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_backups=0
stdout_logfile_backups=0
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
user=nextjs
environment=HOME="/app",USER="nextjs",NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0",LOG_LEVEL="%(ENV_LOG_LEVEL)s"
# Start after PocketBase
priority=200
# Wait for process to stabilize
startsecs=5
startretries=3
# Graceful shutdown: send SIGTERM, wait for connections to drain
stopsignal=TERM
stopwaitsecs=15
killasgroup=true
stopasgroup=true

# =============================================================================
# Background Worker Service
# Priority 300 - Starts after PocketBase is ready
# Requirements: 1.1, 1.2, 13.4
# =============================================================================
[program:worker]
command=node dist/main.js
directory=/app/worker
autostart=true
autorestart=true
# Log to both files and console (stdout/stderr) so errors are visible in Docker logs
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
# Also keep file logs for debugging
stderr_logfile_backups=0
stdout_logfile_backups=0
# Disable log rotation to ensure all output goes to console
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
user=nextjs
environment=HOME="/app",USER="nextjs",NODE_ENV="production",POCKETBASE_URL="%(ENV_POCKETBASE_URL)s",POCKETBASE_ADMIN_EMAIL="%(ENV_POCKETBASE_ADMIN_EMAIL)s",POCKETBASE_ADMIN_PASSWORD="%(ENV_POCKETBASE_ADMIN_PASSWORD)s",WORKER_DATA_DIR="%(ENV_WORKER_DATA_DIR)s",REDIS_URL="%(ENV_REDIS_URL)s",REDIS_HOST="%(ENV_REDIS_HOST)s",REDIS_PORT="%(ENV_REDIS_PORT)s",REDIS_PASSWORD="%(ENV_REDIS_PASSWORD)s",LOG_LEVEL="%(ENV_LOG_LEVEL)s"

# Start after PocketBase and Next.js
priority=300
# Consider started after 5 seconds of running (wait for PocketBase to be ready)
startsecs=5
# Retry up to 5 times on failure (Requirements 1.2: automatic restart)
startretries=5
# Graceful shutdown: send SIGTERM, allow current tasks to complete
stopsignal=TERM
# Wait up to configured timeout for graceful shutdown (Requirements 13.4)
stopwaitsecs=%(ENV_GRACEFUL_SHUTDOWN_TIMEOUT)s
# Ensure all child processes are killed together
killasgroup=true
stopasgroup=true

# =============================================================================
# Nginx Reverse Proxy
# Priority 10 - Starts first to accept connections
# Stops first during shutdown to stop accepting new connections
# =============================================================================
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
# Log to console for visibility
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_backups=0
stdout_logfile_backups=0
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
user=root
# Nginx starts first to accept connections
priority=10
# Graceful shutdown: send SIGQUIT for graceful stop
stopsignal=QUIT
stopwaitsecs=10
killasgroup=true
stopasgroup=true

# =============================================================================
# Supervisor Control Interface
# =============================================================================
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

